# blackfoot/templates/default/gost.service.erb
[Unit]
Description=GOST (Go SensorThings) API service
Wants=postgresql.service

[Service]
Type=simple
User=<%= @user %>
ExecStart=<%= @prefix %>/linux64/gost -config <%= @prefix %>/linux64/config.yaml

[Install]
WantedBy=multi-user.target
